# files/etc/systemd/system/restic@.timer.j2
[Unit]
Description=Run restic@%i.service at specific time

[Timer]
OnCalendar={{ backup_time }}
Unit=restic@%i.service